# debian/python3-bocs.cron.d
################################################################################
#===============================================================================
# Compresses yesterday's data log and removes the raw text file.
################################################################################
0 1 * * * pi bocs_archive 2>&1 | systemd-cat -t bocs_archive
