# etc/systemd/system/bocs_control.service
[Unit]
Description=BOCS Data Collector

[Service]
WorkingDirectory=/home/pi/bocs_control
User=pi
Group=pi
ExecStart=/bin/bash /home/pi/bocs_control/bin/run_control.sh
Type=simple
Restart=always

[Install]
WantedBy=multi-user.target
